# plugins/deprecated/cuda/deploy/cuda-exec.service
[Unit]
Description=cuda_exec GPU kernel compilation, evaluation, and profiling service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.cuda_exec_service
ExecStart=%h/.cuda_exec_service/.venv/bin/python -m uvicorn cuda_exec.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Environment
Environment=PYTHONPATH=%h/.cuda_exec_service
Environment=CUDA_EXEC_ROOT=%h/.cuda_exec
Environment=CUDA_VISIBLE_DEVICES=__CUDA_VISIBLE_DEVICES__
Environment=CUDA_HOME=__CUDA_HOME__
Environment=LD_PRELOAD=__LD_PRELOAD__
Environment=PATH=__CUDA_HOME__/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Hardening directives that imply NoNewPrivileges=true are intentionally
# omitted — NCU profiling requires sudo for GPU performance counters.

# Resource limits
LimitNOFILE=65536
TimeoutStopSec=30

[Install]
WantedBy=default.target
